# deploy/supervisor-laravel-worker.conf
; Supervisor: worker da fila Laravel (jobs como RunAutomationFromActionJob).
; Copie para /etc/supervisor/conf.d/ no servidor e ajuste user e paths.

[program:secretarioapp-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/SecretarioApp/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/SecretarioApp/storage/logs/worker.log
stopwaitsecs=3600
